# etc/s6/dnsmasq/run
#!/bin/sh

exec dnsmasq -C "$DNSMASQ_CONF"
