# etc/s6/tor/run
#!/bin/sh

exec s6-setuidgid tor tor --defaults-torrc "$TOR_CONF"
